# services/backup/cron/backup.service
# =============================================================================
# Systemd Backup Service
# =============================================================================
# Install:
#   sudo cp backup.service /etc/systemd/system/
#   sudo cp backup.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now backup.timer
# =============================================================================

[Unit]
Description=Infrastructure Backup Service
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
ExecStart=/opt/infra/backup/scripts/backup.sh all
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/tmp /opt/infra/backup/logs

[Install]
WantedBy=multi-user.target
